# manyou/install.crontab
0 0 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 0 )
0 1 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 1 )
0 2 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 2 )
0 3 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 3 )
0 4 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 4 )
0 5 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 5 )
0 6 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 6 )
0 7 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 7 )
0 8 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 8 )
0 9 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 9 )
0 10 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 10 )
0 11 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 11 )
0 12 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 12 )
0 13 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 13 )
0 14 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 14 )
0 15 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 15 )
0 16 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 16 )
0 17 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 17 )
0 18 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 18 )
0 19 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 19 )
0 20 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 20 )
0 21 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 21 )
0 22 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 22 )
0 23 * * * ( cd /home/httpd/n030044102/public_html/manyou; /usr/local/bin/python main.py 23 )
